# src/etc/profile.d/Makefile.am
data_SCRIPT=logil-email.sh

EXTRA_DIST=login-email.sh
